# packaging/retroshell.service
[Unit]
Description=RetroShell desktop session (user)
Documentation=file:///usr/share/doc/retroshell/README.md
After=graphical-session-pre.target
Wants=graphical-session-pre.target

[Service]
Type=simple
ExecStart=/usr/local/bin/start-retroshell
Restart=on-failure
RestartSec=2
# Honest fallback still handled inside start-retroshell (labwc if compositor dies).
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=RetroShell

[Install]
WantedBy=default.target
